# files/backup.timer
[Unit]
Description=MythTV backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
